# cpack/debian/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/fluent-bit" ] && [ "$1" = remove ]; then
	invoke-rc.d fluent-bit stop || exit 1
fi
