# cpack/debian/prerm
#!/bin/sh
set -e

case "$1" in
  remove|upgrade|deconfigure)
    if [ -x "/etc/init.d/fluent-bit" ]; then
      invoke-rc.d fluent-bit stop || true
    fi
    ;;
esac
