# lib/nghttp2-1.65.0/contrib/nghttpx-upstart.conf.in
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec @bindir@/nghttpx
